# crontab/root.crontab
# Meant to be installed as root

0 12 */2 * * pacman --noconfirm -Syu >> /var/log/pacman_update.log

# vim:ft=crontab
