# crontab/root.crontab
# Meant to be installed as root

0  12 */2 * * pacman --noconfirm -Syu
20 12 */2 * * pacman -Qdt | xargs pacman --noconfirm -Rns

# vim:ft=crontab
